# supervisord.conf
[program:cored]
command=/opt/core/venv/bin/core-daemon --grpc-address 0.0.0.0 --grpc-port 50051

[program:sshd]
command=/usr/sbin/sshd -D -e

[program:dockerd]
command=/usr/bin/dockerd
